# supervisord.conf
[supervisord]
nodaemon=true

[program:ms_worker]
command=/bin/bash -l -c 'java -jar /app/target/analytics-ms-0.0.1.jar'
process_name=%(program_name)s_%(process_num)02d
autorestart=true
numprocs=1
